# dotfiles/pywalfox/.config/systemd/user/pywalfox-theme-sync.service
[Unit]
Description=Apply the current desktop palette to Firefox

[Service]
Type=oneshot
ExecStart=%h/.local/bin/pywalfox update
