# dotfiles/pywalfox/.config/systemd/user/pywalfox-theme-sync.path
[Unit]
Description=Watch for Pywalfox browser startup

[Path]
PathChanged=/tmp/pywalfox_socket_%U
Unit=pywalfox-theme-sync.service

[Install]
WantedBy=default.target
