# dotfiles/ghostty/.config/systemd/user/app-com.mitchellh.ghostty.service.d/restart.conf
[Service]
# Noctalia may reload its generated Ghostty theme while the service is still
# activating at session startup. Recover from that transient reload race.
Restart=on-abnormal
RestartSec=500ms
